# sv/memcache/run
#!/bin/sh

#sv start socklog-unix || exit 1

exec 2>&1
exec /usr/bin/memcached -u mysql -p 11211 -m 64 -c 1024 -t 4
