# sv/openproject/run
#!/bin/sh

sv start mysql || exit 1

cd /openproject

exec 2>&1
exec bundle exec passenger start -e production 
